# supervisord.conf
[supervisord]
nodaemon=true

[program:tokenpay]
command=/var/www/html/TokenPay   ; 路径替换为你需要执行的脚本或者命令
autostart=true
autorestart=true
redirect_stderr=true
